# conf/supervisor_celery_flower.ini
[program:celery_flower]
command=python manage.py celery flower --address=139.129.92.179 --port=5001
directory=/home/github/mydjango
environment=PATH="/root/anaconda3/envs/mydjango/bin"
stdout_logfile=/home/github/mydjango/logs/celery.flower.log
stderr_logfile=/home/github/mydjango/logs/celery.flower.log
autostart=true
autorestart=true
startsece=10
stopwatisecs=60
priority=1002
