# conf/supervisor_celery_worker.ini
[program:celery_worker]
command=python manage.py celery worker -l INFO
directory=/home/github/mydjango
environment=/root/anaconda3/envs/mydjango/bin
stdout_logfile=/home/github/mydjango/logs/celery.worker.log
stderr_logfile=/hoe/github/mydjango/logs/celery.worker.log
autostart=true
autorestart=true
startsece=10
stopwatisecs=60
priority=100
