# conf/supervisor_celery_beat.ini
[program:celery_beat]
command=python manage.py celery beat -l INFO
directory=/home/github/mydjango
environment=PATH="/root/anaconda3/envs/mydjango/bin"
stdout_logfile=/home/github/mydjango/logs/celery.beat.log
stderr_logfile=/home/github/mydjango/logs/celery.beat.log
autostart=true
autorestart=true
startsece=10
stopwatisecs=60
priority=1001
